# examples/deployment/systemd/tckdb-api.service
# /etc/systemd/system/tckdb-api.service
#
# Runs the TCKDB FastAPI backend on the host (not in Docker) under the
# project's conda env. Used by the Raspberry Pi hosted deployment - see
# docs/deployment/raspberry_pi_hosted_deployment.md.
#
# Install:
#   sudo cp examples/deployment/systemd/tckdb-api.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now tckdb-api.service
#
# Adjust paths/usernames to match the deployment host.

[Unit]
Description=TCKDB FastAPI backend (uvicorn)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=tckdb
Group=tckdb
WorkingDirectory=/opt/tckdb/backend

# Pull env vars (DB creds, abuse-control toggles, etc.) from .env.pi.
# systemd's EnvironmentFile does not perform shell expansion or quoting -
# write `KEY=value` lines exactly as in .env.pi.example.
EnvironmentFile=/opt/tckdb/.env.pi

# Bind to localhost. cloudflared is the only thing that talks to this.
ExecStart=/opt/conda/bin/conda run --no-capture-output -n tckdb_env \
    uvicorn main:app \
        --host 127.0.0.1 \
        --port 8000 \
        --workers 2 \
        --proxy-headers \
        --forwarded-allow-ips 127.0.0.1

Restart=on-failure
RestartSec=5s

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

# Log rotation is handled by journald (see journald.conf) plus the
# Docker json-file driver limits in docker-compose.pi.yml.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
